# kibana/supervisord.conf
[supervisord]
nodaemon=true

[program:http]
command=/usr/bin/twistd -n web --path /opt/kibana
numprocs=1
autostart=true
autorestart=true
